# pi_fireplace.service
[Unit]
Description=Fireplace video looper service

[Install]
WantedBy=multi-user.target

[Service]
ExecStart=/usr/bin/pi_fireplace/start.sh
